# roles/lego/templates/lego@.timer.jinja
# {{ ansible_managed }}

[Unit]
Description=ACME Client Timer for the %I Certificate

[Timer]
OnCalendar=hourly
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target
